# whats-new/zyphor-whats-new/etc/xdg/autostart/zyphor-whats-new.desktop
[Desktop Entry]
Type=Application
Name=Zyphor What's New
Exec=/usr/bin/zyphor-whats-new
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true